# stage3/00-ap-setup/files/ap-setup.service
[Unit]
Description=Access Point Setup Service
After=network.target
Requires=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
  systemctl start hostapd; \
  systemctl start dnsmasq; \
  /usr/bin/python3 /usr/local/bin/setup-server.py & \
  wait'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
