# stage3/01-ap-setup/files/ap-setup.service
[Unit]
Description=Flask Wi-Fi Credential Setup Server
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/setup-server.py
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
