# stage3/01-ap-setup/files/ap-setup.service
[Unit]
Description=Flask AP Setup Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/setup-server.py
Restart=always
StandardOutput=inherit
StandardError=inherit
WorkingDirectory=/usr/local/bin

[Install]
WantedBy=multi-user.target
